# cron/anacrontab
# To set this up, run `anacron -t /home/nsalesky/.dotfiles/cron/anacrontab -S /home/nsalesky/.var/spool/anacron`
SHELL=/usr/bin/fish
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/home/nsalesky/bin
RANDOM_DELAY=30
@daily 5 cron.daily run-parts /home/nsalesky/.dotfiles/cron/daily
@weekly 0 cron.daily run-parts /home/nsalesky/.dotfiles/cron/weekly